# supervisord/supervisord.conf
[supervisord]
nodaemon = true

[program:influxdb]
priority = 1
command = bash -c "/etc/init.d/influxdb start && sleep 5"
autostart=true
autorestart=true

[program:telegraf]
priority = 2
command = bash -c "sleep 15 && /etc/init.d/telegraf start"
autostart=true
autorestart=true

[program:grafana]
priority = 3
command = bash -c "sleep 30 && /etc/init.d/grafana-server start"
autostart=true
autorestart=true
